# deploy/madnesslab.service
[Unit]
Description=MadnessLab FastAPI Backend
After=network.target

[Service]
User=www-data
WorkingDirectory=/opt/madnesslab/backend
ExecStart=/opt/madnesslab/backend/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5
EnvironmentFile=/opt/madnesslab/backend/.env

[Install]
WantedBy=multi-user.target
